# Docker/mms-servers/gitlab-aircon/scripts/crontab
0 0 * * * /media/m2/servers/gitlab/scripts/backup.sh
